# duckdb-mcp-server/deploy/duckdb-mcp-server.service
[Unit]
Description=DuckDB MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/net/ipa.jmcnet/data/h/jmc/src/BERDL-ENIGMA-CORAL/duckdb-mcp-server
Environment=HOST=10.2.2.14
Environment=PORT=80
Environment=DUCKDB_PATH=/scratch/jmc/linkml-coral/cdm_store_bricks_full.db
ExecStartPre=/lab/bin/uv sync --project /mnt/net/ipa.jmcnet/data/h/jmc/src/BERDL-ENIGMA-CORAL/duckdb-mcp-server --python /var/lib/duckdb-mcp/venv/bin/python
ExecStart=/lab/bin/uv run --project /mnt/net/ipa.jmcnet/data/h/jmc/src/BERDL-ENIGMA-CORAL/duckdb-mcp-server --python /var/lib/duckdb-mcp/venv/bin/python src/main.py
Restart=on-failure
RestartSec=3

# Allow binding to privileged port 80 without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# Set these to an appropriate service account.
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
